# rtsp-camera.service
[Unit]
Description=Raspberry Pi RTSP Stream
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/bin/bash %h/rtsp-camera.sh
ExecStartPre=/bin/sleep 10
Restart=always
RestartSec=15
User=CAMERA_USER_PLACEHOLDER
Group=video

# Resource limits for stability
MemoryMax=256M
MemoryAccounting=true
CPUQuota=80%

# Watchdog for detecting hangs
WatchdogSec=300
NotifyAccess=none

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rtsp-stream

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /var/log
PrivateTmp=true

# Environment (systemd will set HOME automatically based on User=)
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
